# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
port=127.0.0.1:9001
username=admin
password=changeme123

[supervisord]
logfile=/Users/arturhanusek/PhpstormProjects/claudeUI/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username=admin
password=changeme123

[program:claudeui]
command=node server.js
directory=/Users/arturhanusek/PhpstormProjects/claudeUI
autostart=true
autorestart=true
stderr_logfile=/Users/arturhanusek/PhpstormProjects/claudeUI/logs/claudeui.err.log
stdout_logfile=/Users/arturhanusek/PhpstormProjects/claudeUI/logs/claudeui.out.log
environment=NODE_ENV="production",PORT="3000"
user=arturhanusek
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true

[eventlistener:claudeui_healthcheck]
command=/Users/arturhanusek/PhpstormProjects/claudeUI/healthcheck.sh
events=TICK_60
stdout_logfile=/Users/arturhanusek/PhpstormProjects/claudeUI/logs/healthcheck.log
stderr_logfile=/Users/arturhanusek/PhpstormProjects/claudeUI/logs/healthcheck.err.log